# multimon400.30.cron
#!/bin/bash
(sleep 1800;cd /var/www/html/multimon;./start.multimon400) &


