# multimon400.cron
#!/bin/bash
(cd /var/www/html/multimon;./start.multimon400) &


